# crontab
0 3 * * 1 /usr/local/bin/python /usr/src/auth_system/project-auth/update_fav_postgres.py >> /tmp/postgres_update.log 2>&1
